# deploy/systemd/njord-md-ingest.service
[Unit]
Description=Njord Market Data Ingest
After=network-online.target

[Service]
User=trade
WorkingDirectory=%h/njord_quant
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/njord_quant/venv/bin/python -m apps.md_ingest --config %h/njord_quant/config/base.yaml
Restart=always
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
